# crontab
*/20 * * * * /usr/bin/mailsync
*/20 * * * * /sbin/castget
